# deploy/truco-bridge.service
[Unit]
Description=Truco Mineiro Online - bridge web (multi-sessão, atrás do Caddy)
After=network.target truco-server.service
Wants=truco-server.service

[Service]
Type=simple
User=truco
WorkingDirectory=/opt/truco-mineiro
# só localhost: quem fala com o mundo de fora agora é o Caddy (porta 443/80
# com HTTPS automático), que repassa pra esse processo internamente.
ExecStart=/usr/bin/python3 /opt/truco-mineiro/client/web_bridge.py 127.0.0.1 5000 8080 127.0.0.1
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
